# .config/systemd/user/swww.timer
[Unit]
Description=Run swww hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
