# .config/systemd/user/walls.timer
[Unit]
Description=Run swww hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
